# cookbooks/nginx/templates/default/dogecoin_gateway.conf.erb
description "dogecoin gateway"
author "Steven Zeiler <steven@zenlabs.co>"
 
start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown
 
respawn
respawn limit 99 5
 
script
export HOME="/home/deploy"
export NODE_ENV=staging"
 
cd /var/www/app/current
exec sudo -u deploy NODE_ENV=staging /usr/bin/node index.js >> /var/log/node.log 2>&1
end script


